# usbip-debian/server/etc_systemd_system_usbipd.service
[Unit]
Description=USB-IP Daemon

[Service]
ExecStart=/usr/sbin/usbipd

[Install]
WantedBy=multi-user.target
